# nullnet_services/i2pd-compose.service
[Unit]
Description=I2PD Compose Service
After=rpm-ostree-install@podman-compose.service

[Service]
Type=oneshot
ExecStart=/usr/bin/podman compose -f /home/nsa_admin/i2pd/compose.yml up -d
ExecStop=/usr/bin/podman compose -f /home/nsa_admin/i2pd/compose.yml stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
